# k2b-board/deploy/app.service
[Unit]
StartLimitIntervalSec=300
StartLimitBurst=5
Description=CrossTalk K2B Client
After=network-online.target pipewire.service
Wants=network-online.target

[Service]
Type=simple
User=streamlate
Group=audio
ExecStart=/usr/local/bin/ct-client
Environment=CROSSTALK_CONFIG=/etc/app/crosstalk.json
Environment=XDG_RUNTIME_DIR=/run/user/999
Environment=PIPEWIRE_RUNTIME_DIR=/run/user/999
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
