# k2b-board/deploy/app.service
[Unit]
StartLimitIntervalSec=300
StartLimitBurst=5
Description=CrossTalk K2B Client
After=network-online.target pipewire.service
Wants=network-online.target

[Service]
Type=simple
User=streamlate
Group=audio
SupplementaryGroups=video
ExecStartPre=+/bin/bash -c '\
  sudo -u streamlate XDG_RUNTIME_DIR=/run/user/999 pactl set-card-profile alsa_card.platform-5096000.codec pro-audio 2>/dev/null || true; \
  sudo -u streamlate XDG_RUNTIME_DIR=/run/user/999 pactl set-card-profile alsa_card.platform-soc_ahub1_mach pro-audio 2>/dev/null || true; \
  amixer -c 1 cset numid=4 on 2>/dev/null || true; \
  amixer -c 1 cset numid=5 on 2>/dev/null || true; \
  amixer -c 1 cset numid=6 on 2>/dev/null || true; \
  amixer -c 1 cset numid=7 on 2>/dev/null || true; \
  amixer -c 1 cset numid=2 31 2>/dev/null || true; \
  amixer -c 1 cset numid=3 on 2>/dev/null || true'
ExecStart=/usr/local/bin/ct-client
Environment=CROSSTALK_CONFIG=/etc/app/crosstalk.json
Environment=XDG_RUNTIME_DIR=/run/user/999
Environment=PIPEWIRE_RUNTIME_DIR=/run/user/999
Environment=PULSE_RUNTIME_PATH=/run/user/999/pulse
Environment=USE_DISPLAY=true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
